# systemd/dbmail-timsieved.service.in
[Unit]
Description=DBMail Sieve Server
After=network.target network-online.target mysqld.service postgresql.service

[Service]
Type=notify
ExecStart=@sbindir@/dbmail-timsieved -D
Restart=always
RestartSec=3s

[Install]
WantedBy=multi-user.target
