# systemd/dbmail-lmtpd.service.in
[Unit]
Description=DBMail LMTP Server
After=syslog.target network.target mysqld.service postgresql.service

[Service]
Type=forking
PIDFile=@piddir@/dbmail-lmtpd.pid
ExecStart=@sbindir@/dbmail-lmtpd -p @piddir@/dbmail-lmtpd.pid

[Install]
WantedBy=multi-user.target
