# systemd/dbmail-imapd.service.in
[Unit]
Description=DBMail Imap Server
After=syslog.target network.target mysqld.service postgresql.service

[Service]
Type=notify
ExecStart=@sbindir@/dbmail-imapd -D

[Install]
WantedBy=multi-user.target
