# systemd/dbmail-lmtpd.service.in
[Unit]
Description=DBMail LMTP Server
After=network.target network-online.target mysqld.service postgresql.service

[Service]
Type=notify
ExecStart=@sbindir@/dbmail-lmtpd -D
Restart=always
RestartSec=3s

[Install]
WantedBy=multi-user.target
