# systemd/dbmail-pop3d.service.in
[Unit]
Description=DBMail pop3 Server
After=network.target mysqld.service postgresql.service

[Service]
Type=simple
ExecStart=@sbindir@/dbmail-pop3d -D
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
