# debian/dbmail.cron.d
0 3 * * * root /usr/sbin/dbmail-util -c -t -u -r -p -d &>/dev/null
