# debian/dbmail.cron.d
0 3 * * * root /usr/sbin/dbmail-util -ay >/dev/null 2>&1
