# debian/dbmail.cron.d
0 3 * * * root text -x /usr/sbin/dbmail-util && /usr/sbin/dbmail-util -ay >/dev/null 2>&1
